# wagtail-supervisor.conf
[program:wagtail]
command=/wagtaildemo/venv/bin/gunicorn wagtaildemo.wsgi:application
directory=/wagtaildemo
user=www-data
autostart=true
autorestart=true
